# etc/service/dbus-daemon/run
#!/bin/sh
rm -f /var/run/dbus/*
exec dbus-daemon --nofork --system