# etc/service/dbus-daemon/run
#!/bin/sh
exec dbus-daemon --nofork --system